# pwmfan.service
[Unit]
Description=Raspberry Pi PWM Fan Controller
Documentation=https://github.com/yourusername/pypwm
After=network.target
Wants=network.target
After=raspberrypi-sys-mods.service

[Service]
Type=simple
# Running as root for GPIO access
User=root
Group=root
ExecStart=/home/amirulandalib/Scripts/pypwm/venv/bin/python /home/amirulandalib/Scripts/pypwm/daemon.py
WorkingDirectory=/home/amirulandalib/Scripts/pypwm

# Environment settings
Environment=PYTHONUNBUFFERED=1
Environment=BASE_DIR=/home/amirulandalib/Scripts/pypwm

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
# Modified security settings for root access
ProtectSystem=false
PrivateTmp=false
NoNewPrivileges=false

# Resource limits
CPUQuota=25%
MemoryLimit=100M

# Logging
StandardOutput=append:/home/amirulandalib/Scripts/pypwm/logs/daemon.log
StandardError=append:/home/amirulandalib/Scripts/pypwm/logs/daemon.error.log

[Install]
WantedBy=multi-user.target